# go-peer/aleph/rootfs/uc-go-peer.service
[Unit]
Description=universal-connectivity go-peer relay node
Wants=network-online.target
After=network-online.target
ConditionPathExists=/etc/default/uc-go-peer.ready

[Service]
Type=simple
User=uc-go-peer
Group=uc-go-peer
WorkingDirectory=/var/lib/uc-go-peer
EnvironmentFile=-/etc/default/uc-go-peer
ExecStart=/bin/sh -lc 'exec /usr/local/bin/universal-chat-go --headless -identity "${GO_PEER_IDENTITY_PATH:-/var/lib/uc-go-peer/identity.key}" -port-tcp "${GO_PEER_TCP_PORT:-9095}" -port-ws "${GO_PEER_WS_PORT:-9096}" -port-quic "${GO_PEER_QUIC_PORT:-9097}" -port-webtransport "${GO_PEER_WEBTRANSPORT_PORT:-9097}" -port-webrtc "${GO_PEER_WEBRTC_PORT:-9098}"'
Restart=always
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
